# files/host-veg/xdg_config_home/systemd/user/firefox.slice
[Unit]
Description=Firefox Web Browser

[Install]
WantedBy=user-%U.target

[Slice]
#CPUWeight=256
CPUWeight=400
#CPUQuota=20%
MemoryMax=8G
