# files/host-veg/xdg_config_home/systemd/user/polkit-gnome-authentication-agent@.service
[Unit]
Description=polkit-GNOME authentication agent
BindsTo=xsession-start@%i.service

[Service]
EnvironmentFile=%h/.xsession.env
PassEnvironment=DESKTOP_SESSION
Type=simple
ExecStart=/usr/libexec/polkit-gnome-authentication-agent-1
